# supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:updater]
command=/opt/copy-file.sh
autostart=true

[program:snort]
command=/opt/entrypoint.sh
autostart=true
